# supervisord.d/nginx.ini
[program:nginx]
command=/usr/local/nginx/sbin/nginx -g 'daemon off'
startsecs=0
autostart=true
autorestart=true
stdout_logfile=/usr/local/nginx/logs/error.log
stopasgroup=true
killasgroup=true
